# rootfs/etc/s6/services/20-postfix/run
#!/usr/bin/with-contenv sh
sleep 2
echo "[*] Starting postfix"
exec /usr/sbin/postfix start-fg
