# rootfs/etc/s6/services/30-dovecot/run
#!/usr/bin/with-contenv sh
sleep 3
echo "[*] Starting dovecot"
exec /usr/sbin/dovecot -F
