# rootfs/etc/s6/services/10-rsyslogd/run
#!/usr/bin/with-contenv sh
echo "[NOTICE] Starting rsyslog"
exec /usr/sbin/rsyslogd -n
